# xinetd.d/xinetd.d/httpd-7
service http-7
{
  disable = no
  type = UNLISTED
  socket_type = stream
  protocol = tcp
  wait = no
  redirect = 192.168.0.18 80
  bind = 0.0.0.0
  port = 7080
  user = nobody
}
